# mac-mini-m4/launchd/com.local.z2m-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.z2m-backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alex/komodo-dean-gitops/mac-mini-m4/scripts/z2m-backup.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>UserName</key>
    <string>alex</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/alex</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/z2m-backup.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/z2m-backup.log</string>
</dict>
</plist>
