# launchd/com.local.komodo-stack-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.komodo-stack-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alex/mac-mini-compose/scripts/sync-stacks.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/alex</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/komodo-stack-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/komodo-stack-sync.log</string>
</dict>
</plist>
